# recipes-fabmo/fabmo-updater/files/factory-reset.service
[Unit]
Description=Factory Reset

[Service]
Type=oneshot
ExecStart=/usr/lib/fabmo/factory_reset.sh

[Install]
WantedBy=default.target

